# Update Client/systemd/dyndns-auto-update.timer
[Unit]
Description=Starts the dyndns auto update client

[Timer]
OnCalendar=*-*-* *:0/5:00
Persistent=true

[Install]
WantedBy=timers.target
